# contrib/systemd/user/waytrim-watch.service
[Unit]
Description=waytrim automatic clipboard cleaner
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/waytrim-watch auto
Restart=on-failure
RestartSec=1

[Install]
WantedBy=default.target
